# deploy/telegram-forwarder-webhook.service
[Unit]
Description=Telegram Forwarder Razorpay Webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/telegram-forwarder-bot
EnvironmentFile=/opt/telegram-forwarder-bot/.env
ExecStart=/opt/telegram-forwarder-bot/.venv/bin/uvicorn razorpay_webhook:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
